# default-lldpad/templates/etc/cron.d/lldpad.j2
# Run lldpad script every 10 minutes
*/10 * * * * root "{{ lldp_path_file }}"
